# packaging/opentelemetry-injector/etc/ld.so.preload
/opt/opentelemetry/libopentelemetry_injector.so